# docker/workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/console schedule:run >> /dev/null 2>&1
